# scripts/cron/rcas6004-starqa.cron
    */10 * * * * /direct/star+u/starqa/QA/daemon > /dev/null 2> /dev/null


