# scripts/cron/rcas6004-starqa.cron
* * * * * /direct/star+u/starqa/bin/QAdaemon > /dev/null 2> /dev/null
0 * * * * /star/u/starqa/QA/newFilesCheck_NoDB > /dev/null 2> /dev/null
30 * * * * /star/u/starqa/QA/newFilesCheck_NoDB > /dev/null 2> /dev/null

