# scripts/cron/rplay17-starlib.cron
00 05 * * *     /afs/rhic.bnl.gov/star/packages/scripts/AutoBuild.csh SL302
